# systemd/livekit-wake.service
[Unit]
Description=LiveKit Wake Word Listener
After=livekit-agent.service sound.target

[Service]
User=pi
WorkingDirectory=/home/pi/jarvis-assistant/agent
EnvironmentFile=/home/pi/jarvis-assistant/.env
ExecStart=/home/pi/jarvis-assistant/.venv/bin/python main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
