# systemd/livekit-agent.service
[Unit]
Description=LiveKit Voice Agent Worker
After=network.target docker.service
Requires=docker.service

[Service]
User=pi
WorkingDirectory=/home/pi/pi-voice-agent/agent
EnvironmentFile=/home/pi/pi-voice-agent/.env
ExecStart=/home/pi/pi-voice-agent/.venv/bin/python agent.py start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
